# thin_edge_example_python_module/debian/python3-thin-edge-example-python-module.service
[Unit]
Description=This service is an example for a python module for thin-edge

[Service]
Type=simple
ExecStart=/usr/bin/thin-edge-module-example
Restart=always

[Install]
WantedBy=multi-user.target