# thin_edge_file_management/debian/thin-edge-file-management-python3.service
[Unit]
Description=A thin-edge.io plugin to manage configuration and log files

[Service]
Type=simple
ExecStart=/usr/bin/thin-edge-file-management
Restart=always

[Install]
WantedBy=multi-user.target