# thin_edge_tf_linear_poti/debian/thin-edge-tf-linear-poti-python3.service
[Unit]
Description=A thin-edge.io plugin to integrate TinkerForge Linear Poti

[Service]
Type=simple
ExecStart=/usr/bin/thin-edge-tf-linear-poti
Restart=always

[Install]
WantedBy=multi-user.target